# app.socket
[Unit]
Description=app socket

[Socket]
ListenStream=/run/app.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
